# deploy/railway/supervisord.conf
; supervisord.conf — runs as PID 1 inside the Railway container
; Manages: gunicorn (web), node socketio, 2 RQ workers, scheduler, nginx

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ---------- Frappe gunicorn web ----------
[program:frappe-web]
command=/home/frappe/frappe-bench/env/bin/gunicorn
    -b 127.0.0.1:8000
    --workers 4
    --timeout 120
    --worker-tmp-dir /dev/shm
    frappe.app:application
directory=/home/frappe/frappe-bench/sites
user=frappe
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---------- Frappe socketio (realtime) ----------
[program:frappe-socketio]
command=/usr/bin/node /home/frappe/frappe-bench/apps/frappe/socketio.js
directory=/home/frappe/frappe-bench
user=frappe
environment=NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---------- RQ worker: short queue ----------
[program:frappe-worker-short]
command=/home/frappe/.local/bin/bench worker --queue short
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe",PATH="/home/frappe/.local/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---------- RQ worker: long + default queues ----------
[program:frappe-worker-long]
command=/home/frappe/.local/bin/bench worker --queue long,default
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe",PATH="/home/frappe/.local/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---------- Scheduler ----------
[program:frappe-schedule]
command=/home/frappe/.local/bin/bench schedule
directory=/home/frappe/frappe-bench
user=frappe
environment=HOME="/home/frappe",PATH="/home/frappe/.local/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---------- Nginx (HTTP entrypoint) ----------
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
